# src/main/resources/systemd-template.service
[Unit]
Description=Minecraft Server
After=network.target

[Install]
WantedBy = default.target

[Service]
Restart = always
